# supervisord.conf
[program:node]
command = node bin/trellocalendar -p 8080
directory = /home/dotcloud/current
