# supervisord.conf
[program:node]
command = node bin/trellocalendar
directory = /home/dotcloud/current
